# deploy/systemd/ester.service
# deploy/systemd/ester.service — unit-fayl systemd
# MOSTY:
# - (Explicit) Launches uvikorn asgi.app_mine:app under the given user.
# - (Hidden #1) Reads variables from /ets/ester.env (create a copy of .env).
# - (Hidden #2) Restart=always + simple health-logic via systemd.

[Unit]
Description=Ester Synergy Service
After=network.target

[Service]
Type=simple
User=%i
EnvironmentFile=/etc/ester.env
WorkingDirectory=/opt/ester
ExecStart=/usr/bin/python -m uvicorn asgi.app_main:app --host ${HOST:-0.0.0.0} --port ${PORT:-8080} --workers ${WORKERS:-2} --log-level ${LOG_LEVEL:-info}
Restart=always
RestartSec=5

# Zhurnaly
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
