# systemd/ester-crdt-snapshot.service
[Unit]
Description=Ester CRDT Snapshot Export (CAS-backed)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/ester
Environment=PYTHONUNBUFFERED=1
# Ukazhi peremennye pri neobkhodimosti JWT/portov:
# Environment=JWT=<token_with_replicator>
# Environment=PORT=5000
ExecStart=/bin/sh -c 'curl -sSf -X POST ${JWT:+-H "Authorization: Bearer ${JWT}"} http://127.0.0.1:${PORT:-5000}/p2p/snapshot/export > %h/ester/dist/crdt_snapshot_$(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ).json'
