# systemd/ester-quorum-engine.service
# systemd/ester-quorum-engine.service — demon kvoruma (2-of-3)
[Unit]
Description=Ester Quorum Engine (2-of-3 approvals)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ester
Group=ester
WorkingDirectory=/opt/ester
Environment=PYTHONUNBUFFERED=1
# Environment variables can be placed in /ets/ester/ester.env
EnvironmentFile=-/etc/ester/ester.env
# The module is launched as a package: pothon -m modules.approvals.corum_engine
ExecStart=/usr/bin/python3 -m modules.approvals.quorum_engine
Restart=on-failure
RestartSec=2s
# Security restrictions (enable if possible)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/ester /var/log/ester /var/ester /var/quorum
# Logs to syslog/log
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
