# systemd/ester-judge-overlay-watcher.service
# packaging/systemd/ester-judge-overlay-watcher.service
# Yuje overlay file guard unit.

# Mosty:
# - Yavnyy (Ekspluatatsiya ↔ Orkestratsiya): podderzhivaet “zhivoy” kontur primeneniya podskazok.
# - Skrytyy 1 (Bezopasnost ↔ Praktika): komanda perezagruzki zadaetsya snaruzhi (ENV) i mozhet byt empty.
# - Skrytyy 2 (Infoteoriya ↔ Nadezhnost): prostaya pull-model — menshe tochek otkaza.

# Zemnoy abzats:
# The daemon reads ~/.ester/yuje_slots.jsion and when changed, calls the Yuje restart command (if given).

[Unit]
Description=Ester Judge Overlay Watcher
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/ester-judge-overlay-watcher
ExecStart=/usr/bin/env python3 -m scripts.judge_overlay_watch
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
