# systemd/ester-auto-backup.service
[Unit]
Description=Ester — push latest archives to configured storage targets
Documentation=man:systemd.service(5)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# The path to the project directory can be redefined via Environment=ESTER_DIR=/path/to/ester
Environment=ESTER_DIR=/opt/ester
# You can also override the data directory: Environment=PERSIST_DIR=/var/lib/ester
# No dangerous privileges are required; if necessary, set the user/group:
# User=ester
# Group=ester

# We execute the archive auto-publishing script.
# By default, it reads targets from modules.storage.targets and archives from self/archive.
ExecStart=/usr/bin/env bash -lc 'cd "${ESTER_DIR:-/opt/ester}" && PERSIST_DIR="${PERSIST_DIR:-${ESTER_DIR:-/opt/ester}/data}" python3 scripts/auto_backup.py'

# Convenient log output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ester-auto-backup

# No restarts are needed for oneshut
Restart=no

[Install]
WantedBy=multi-user.target
